# deployment/myapp.conf
[program:myapp]
command=/home/ubuntu/myapp/venv/bin/gunicorn -b localhost:8000 -w 4 myapp:app
directory=/home/ubuntu/myapp
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true